# apply_tc_rules.cron
# Applies tc traffic limiting/filtering rules for the ndt_e2e.sh test.

MAILTO=root

@daily root /bin/apply_tc_rules.sh
